# SCAPE604-node-aireas.conf
# SCAPE604-node-aireas - start or restart node-aireas service
#
#

description     "(re)start node-aireas service"
start on runlevel [2345]
respawn
task
script
   /opt/SCAPE604/aireas-data-pg/node-aireas.sh /opt/SCAPE604/log/SCAPE604-node-aireas.log
end script